# .history/gunicorn_20221124025245.service
[Unit]
Description=gunicorn daemon
Requires=gunicorn.socket
After=network.target

[Service]
User=root
Group=www-data
WorkingDirectory=/var/lib/jenkins/workspace/Empact_api/e
ExecStart=/var/lib/jenkins/workspace/Empact_api/env/bin/gunicorn empact.wsgi:application

[Install]
WantedBy=multi-user.target
